# dist/comms.service
[Unit]
Description=comms messaging daemon
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=%h/go/bin/comms daemon start
WorkingDirectory=%h
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
